# Lab2/poweroff.service
[Unit]
Description=Friday 22:00 system shutdown

[Service]
Type=oneshot
ExecStart=/home/milax/GSX/GestionSistemas/Lab2/poweroff.sh
